# script/cron/clubcar_deactivate.cron
#!/bin/sh

dir=`dirname $0`

exec $dir/../runner 'User.find_all_by_present(true).each {|u| u.present = false ; u.save }'
