# script/cron/clubcar_choices.cron
#!/bin/sh

dir=`dirname $0`

exec $dir/../runner 'require RAILS_ROOT + "/app/controllers/application.rb"; RestaurantsController.new.mail_choices'
