# script/wrappers/restart.cron
#!/bin/sh

. /afs/slackers.net/projects/infrastructure/etc/functions
getToken webmaster webmaster 10m

scriptdir=`dirname $0`
basedir=`readlink -f $scriptdir/../release/current`

cd $basedir

pid=`cat tmp/pids/dispatch.3000.pid`
kill -TERM $pid

while kill -0 $pid 2>/dev/null ; do
  sleep 1
done

nohup ./script/spin >/dev/null 2>&1